# deploy/systemd/luarca-ais-web.service
[Unit]
Description=Luarca AIS static web server
After=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__HOME__/web
ExecStart=__HOME__/.venv/bin/python -m http.server __PORT__ --bind 127.0.0.1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Sandboxing
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true

[Install]
WantedBy=multi-user.target
